# supervisord/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
loglevel=info

[program:octane]
command=/usr/bin/env php -d variables_order=EGPCS /var/www/artisan octane:start --watch --server=swoole --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[supervisorctl]
nodaemon=true